# php-fpm.ini
[program:php]
command = /opt/remi/php70/root/usr/sbin/php-fpm
user = nginx
autostart = true